# shelley.service
[Unit]
Description=Shelley-the-Agent
BindsTo=shelley.socket
After=shelley.socket

[Service]
Type=exec
User=exedev
Group=exedev
WorkingDirectory=/home/exedev
ExecStart=/usr/local/bin/shelley -debug -db /home/exedev/.config/shelley/shelley.db -config /exe.dev/shelley.json serve -systemd-activation -require-header X-Exedev-Userid
Restart=on-failure
RestartSec=5
Environment=HOME=/home/exedev
Environment=USER=exedev
Environment=PATH=/headless-shell:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
KillMode=process
TimeoutStopSec=10

StandardOutput=journal
StandardError=journal
SyslogIdentifier=shelley

[Install]
WantedBy=multi-user.target
